# tag-osx/Library/LaunchAgents/remap-keys.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>hu.lackac.remap-keys</string>
        <key>ProgramArguments</key>
        <array>
          <string>sudo</string>
          <string>/Users/lackac/bin/remap_keys</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>KeepAlive</key>
        <false/>
    </dict>
</plist>
